# install/php-fpmadmin/supervisor/system-redis.conf
[program: mysite]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --sleep=3 --tries=1
directory=/var/www/html/mysite
autostart=true
autorestart=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/html/mysite/storage/logs/queue.log
